# deploy/qq-bot.service
[Unit]
Description=QQ Bot Bridge Service
After=network-online.target napcat.service
Wants=network-online.target napcat.service
BindsTo=napcat.service

[Service]
Type=simple
User=roux
WorkingDirectory=/home/roux/Claude_code_test
Environment="PATH=/usr/bin:/usr/local/bin:/home/roux/.local/bin"
ExecStart=/usr/bin/python3 -m qq_bot
Restart=always
RestartSec=15
StandardOutput=append:/home/roux/bridge.log
StandardError=append:/home/roux/bridge.log

# 防止无限重启耗尽资源
StartLimitInterval=300
StartLimitBurst=10

[Install]
WantedBy=multi-user.target
